# nginx/srcs/supervisord.conf
[supervisord]
pidfile=/run/supervisord.pid    ; supervisord pidfile; default supervisord.pid
nodaemon=true                  ; start in foreground if true; default false

[program:nginx]
command=nginx -g "daemon off;"
