# srcs/InfluxDB/srcs/supervisord.conf
[supervisord]
pidfile=/run/supervisord.pid    ; supervisord pidfile; default supervisord.pid
nodaemon=true                  ; start in foreground if true; default false

[program:influxdb]
command=influxd -config /etc/influxdb.conf

[program:telegraf]
command=telegraf -config  /etc/telegraf.conf

[eventlistener:processes]
command=/tmp/stop.sh
events=PROCESS_STATE_STOPPED, PROCESS_STATE_EXITED, PROCESS_STATE_FATAL